# builds/DEBIAN/postinst
#!/bin/sh
set -e
systemctl enable usbipd
systemctl enable qusbipserver
deb-systemd-invoke start usbipd
deb-systemd-invoke start qusbipserver
